# docker/haproxy/etc/cron.daily/letsencrypt
#!/bin/sh

cd /usr/local/etc/haproxy
cat fullchain.pem privkey.pem > cert.pem
#